# exozyme/lydia.service
[Unit]
Description=lydia server

[Service]
ExecStart=~/lydia/d/run.sh

[Install]
WantedBy=default.target
