# exozyme/lydia-update.timer
[Unit]
Description=Update lydia deployment

[Timer]
OnCalendar=*:0..59/5
Persistent=true

[Install]
WantedBy=timers.target