# services/mediawiki/crontab
0 4 * * * root <% $(realpath "$STATES_DIR") %>/atre svc mediawiki cronjob
