# services/mediawiki/crontab
0 4 * * tue,fri root <% $(realpath "$STATES_DIR") %>/services/mediawiki/script/cronjob.sh
