# services/mediawiki/crontab
0 4 * * tue,fri root <% $(realpath "$STATES_DIR") %>/atre svc mediawiki cronjob
