# toolforge/crontab
*/15 * * * * /usr/bin/jsub -N update -once -quiet ./d/cron.sh