# toolforge/crontab
*/15 * * * * ./d/cron.sh >update.log 2>&1