# services/atremis/systemd/atre-pull.timer
[Unit]
Description=Pull Atremis configurations automatically

[Timer]
OnStartupSec=30s
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
