# services/atremis/systemd/atre-pull.timer
[Unit]
Description=Pull Atremis configurations automatically

[Timer]
OnStartupSec=5
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
